# basic/target/linux/deb/postrm
#!/bin/sh

update-rc.d -f muniubox-platform remove
