# basic/target/linux/muniubox-platform.service
[Unit]
Description=muniubox-platform
After=syslog.target

[Service]
User=muniubox
ExecStart=/usr/share/muniubox-platform/bin/muniubox-platform.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target